# truthdb.service
[Unit]
Description=TruthDB service
Wants=network-online.target
After=network-online.target

[Service]
Type=simple
User=truthdb
Group=truthdb

ExecStart=/usr/local/bin/truthdb
WorkingDirectory=/var/lib/truthdb

Restart=on-failure
RestartSec=2
TimeoutStopSec=30
KillSignal=SIGTERM

StateDirectory=truthdb
RuntimeDirectory=truthdb

NoNewPrivileges=yes
PrivateTmp=yes
ProtectSystem=strict
ProtectHome=true
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
LockPersonality=yes
RestrictRealtime=yes
RestrictSUIDSGID=yes
MemoryDenyWriteExecute=yes
UMask=0077

ReadWritePaths=/var/lib/truthdb /run/truthdb

[Install]
WantedBy=multi-user.target
